# deploy/config/containerd-pid-cleanup.timer
[Unit]
Description=Hourly cleanup of stale containerd exec PID files

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
