# deploy/config/togather-docker-prune.timer
[Unit]
Description=Weekly Docker prune for Togather

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
